# SystemD/.config/systemd/user/backup.timer
[Unit]
Description=Run backup weekly

[Timer]
OnUnitActiveSec=1w
Unit=backup.service

[Install]
WantedBy=timers.target
