# SystemD/.config/systemd/user/backup.service
[Unit]
Description=Taking backup

[Service]
Type=simple
ExecStart=%h/.local/scripts/do_backup 
