# SystemD/.config/systemd/user/backup.service
[Unit]
Description=Taking backup

[Service]
Type=simple
ExecStart=/usr/bin/rsync -a --delete --quiet /home/oreo /media/hdd/backup
