# SystemD/.config/systemd/user/mail.service
[Unit]
Description=Syncing mail using mbsyncrc

[Service]
Type=simple
ExecStart=%h/.local/scripts/mail_fetch
