# SystemD/.config/systemd/user/mail.timer
[Unit]
Description=Run mail sync every 3 hours

[Timer]
OnBootSec=1min
OnUnitActiveSec=10800
Unit=mail.service

[Install]
WantedBy=timers.target
