# SystemD/.config/systemd/user/mail.timer
[Unit]
Description=Run mail sync every hour

[Timer]
OnBootSec=1min
OnUnitActiveSec=1800
Unit=mail.service

[Install]
WantedBy=timers.target
